# crontabs/zabbix_autoupdate
# Pull latest version from git.
15 0 * * * root cd /opt/zabbixutil/; git pull -q && /opt/zabbixutil/packages/debian/install
# 
