# packages/debian/ssh.authorizedkeys/crontab
#!/bin/bash
# Cronjob for keeping track of security updates

11 */4 * * * root /usr/local/bin/zabbix.debian.ssh.authorizedkeys > /var/run/zabbix/debian.ssh.authorizedkeys
#
